# private-cookbooks/pxe-server/templates/default/rc.local.erb
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# kitsunai Ubuntuのディスクイメージを読み込み専用で自動マウント
mount -o ro /var/lib/tftpboot/ubuntu-13.10-desktop-amd64.iso /var/lib/tftpboot/image/ubuntu-13.10-desktop-amd64
# kitsunai Ubuntuのディスクイメージを読み込み専用で自動マウント end

exit 0
